# docker/rootfs/etc/systemd/system/bootstrap.service
[Service]
Type=oneshot
ExecStart=/usr/local/bin/bootstrap.sh

[Install]
WantedBy=multi-user.target
