# docker/rootfs/etc/systemd/system/passenger-exporter.service
[Service]
ExecStart=/usr/local/bin/passenger-exporter

[Install]
WantedBy=multi-user.target
