# rootfs/etc/systemd/system/fix-webapp-permissions.service
[Service]
Type=oneshot
ExecStart=/usr/local/bin/fix-webapp-permissions.sh

[Install]
WantedBy=multi-user.target
