# rootfs/etc/systemd/system/root-boot.service
[Service]
Type=oneshot
ExecStart=/usr/local/bin/root_boot.sh

[Install]
WantedBy=multi-user.target
